# systemd/site-monitor.service
[Unit]
Description=Site monitoring check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=MONITOR_SOFT_FAIL=1
ExecStart=/usr/bin/flock -n -E 99 /tmp/site-monitor.lock /root/training/site/run-sites.sh


[Install]
WantedBy=multi-user.target
